# deploy/hot-reload/dcc-hot-reload.service
[Unit]
Description=DCC hot reload Discord Gateway host
Documentation=file:/usr/local/share/dcc/deploy/hot-reload/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dcc
Group=dcc
WorkingDirectory=/var/lib/dcc/hot_reload
EnvironmentFile=/etc/dcc/dcc-hot-reload.env
ExecStartPre=/usr/local/bin/dcc_hot_reload_host --check
ExecStart=/usr/local/bin/dcc_hot_reload_host
Restart=on-failure
RestartSec=2s
KillSignal=SIGTERM
TimeoutStopSec=30s

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/dcc /var/log/dcc
CapabilityBoundingSet=
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
